# ligaLampProjDiebold/DEBIAN/postinst
#!/bin/bash -e
update-rc.d poff.sh defaults
update-rc.d pon.sh defaults
